# deploy/refinery.service
[Unit]
Description=Refinery Balance
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/refinery-balance/backend
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
